# config/task.cron
# m h  dom mon dow   command
0 0 * * * /bin/logrotate.sh >> /var/log/cron.log 2>&1
